# etc/profile.d/lang.sh
export LANG=en_US.UTF-8
export LC_MONETARY=es_PE.UTF-8
export LC_TIME=es_PE.UTF-8
export LC_MEASUREMENT=es_PE.UTF-8
export LC_COLLATE=C
